# cron.d/myexportmail
0 2 * * * root /usr/bin/myexportmail -c /etc/myexportmail/myexportmail.conf > /dev/null 2>&1
